# ansible/roles/teleport-get-cert/templates/teleport-get-cert.service.j2
[Unit]
Description=Service getting teleport certificates

[Service]
Type=oneshot
ExecStartPre={{ teleport_bin_dir }}teleport-lock
ExecStart={{ teleport_bin_dir }}teleport-get-cert

