# ansible/roles/teleport-get-cert/templates/teleport-get-cert.timer.j2
[Unit]
Description=Timer getting letsencrypt certificates

[Timer]
OnBootSec=5min
OnCalendar=hourly
Persistent=true

