# ansible/roles/teleport-push-cert/templates/teleport-push-cert.service.j2
[Unit]
Description=Service pushing teleport certificates from S3

[Service]
Type=oneshot
ExecStart={{ teleport_bin_dir }}teleport-push-cert

