# ansible/roles/teleport-push-cert/templates/teleport-push-cert.timer.j2
[Unit]
Description=Timer pushing letsencrypt certificates from S3

[Timer]
OnBootSec=5min
OnCalendar=hourly
Persistent=true
