# ansible/roles/teleport-ssm-publish-tokens/templates/teleport-ssm-publish-tokens.timer.j2
[Unit]
Description=Timer rotating teleport tokens in SSM

[Timer]
OnBootSec=5min
OnCalendar=hourly
Persistent=true
